# dot_config/systemd/user/codestats-relay.service
[Unit]
Description=Relay HTTP service for codestats.net which returns immidiately to prevent lags in editors

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'codestats-relay'
Restart=on-failure

[Install]
WantedBy=default.target
